# raspberry/systemd/shroom-mqtt-bridge.service
[Unit]
Description=Shroom MQTT to Sensors Ingest Bridge
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Type=simple
User=__RUN_USER__
Group=__RUN_USER__
WorkingDirectory=__WORKDIR__
EnvironmentFile=__ENV_FILE__
ExecStart=__PYTHON_BIN__ __SCRIPT_PATH__
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
